# systemd/icy-flac-proxy.service
[Unit]
Description=Radio ICY-metadata proxy for the FLAC/Ogg mount
After=network-online.target icecast2.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/radio
EnvironmentFile=-/etc/musicbestman/env
ExecStart=/opt/radio/venv/bin/python3 /opt/radio/icy_flac_proxy.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
