# systemd/publish-listeners.service
[Unit]
Description=Radio Listener Count Publisher
After=network-online.target icecast2.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/radio
EnvironmentFile=-/etc/musicbestman/env
ExecStart=/opt/radio/venv/bin/python /opt/radio/publish_listeners.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
